# systemd/user/synthia-frontend-dev.service.in
[Unit]
Description=Synthia Frontend Dev Service (Vite on 5173)
After=network.target

[Service]
Type=simple
WorkingDirectory=@INSTALL_DIR@/frontend
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
