# systemd/user/synthia-supervisor.service.in
[Unit]
Description=Synthia Supervisor (Standalone Addon Orchestrator)
After=network.target

[Service]
Type=simple
WorkingDirectory=@INSTALL_DIR@/backend
EnvironmentFile=-@INSTALL_DIR@/scripts/synthia.env

# Where standalone services live (sibling folder to SynthiaCore)
Environment=SYNTHIA_ADDONS_DIR=%h/Projects/SynthiaAddons
Environment=SYNTHIA_CATALOG_PUBLISHERS=@INSTALL_DIR@/runtime/store/cache/official/publishers.json
Environment=SYNTHIA_SUPERVISOR_INTERVAL_S=5
Environment=PYTHONUNBUFFERED=1

ExecStart=@INSTALL_DIR@/backend/.venv/bin/python -m synthia_supervisor.main

Restart=always
RestartSec=2
EnvironmentFile=%h/.config/synthia/admin.env

[Install]
WantedBy=default.target
